# ansible/idebackend.service
[Service]
Restart=on-failure
WorkingDirectory=/home/ec2-user/ide-backend
ExecStart=/usr/bin/npm run start
